# macfanctld.service
[Unit]
Description=Macfand

[Service]
ExecStart=/usr/sbin/macfanctld

[Install]
WantedBy=multi-user.target
